# infra/monitoring/systemd/post_eval_watcher.service
[Unit]
Description=Nexus Alpha - post_eval_watcher
After=network.target

[Service]
User=azt
WorkingDirectory=/Users/azt/Desktop/Python/nexus-alpha/infra/self_healing
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 post_eval_watcher.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
